# deploy/bft.service
[Unit]
Description=Blue Force Tracker
After=network.target

[Service]
Type=simple
User=jhh-pi
WorkingDirectory=/opt/blufor-tracker
ExecStart=/opt/blufor-tracker/deploy/start.sh
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
