# deploy/fts-ui.service
[Unit]
Description=FreeTAKServer Web UI
After=network.target fts.service
Wants=fts.service

[Service]
Type=simple
User=jhh-pi
ExecStart=/opt/blufor-tracker/.fts-venv/bin/python -c "from FreeTAKServer_UI import run; run.app.run(host='0.0.0.0', port=5000)"
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
